# res/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file
username=aaa                ; default is no username (open server)
password=acell              ; default is no password (open server)

[inet_http_server]         ; inet (TCP) server disabled by default
port=0.0.0.0:23231         ; ip_address:port specifier, *:port for all iface
username=aaa               ; default is no username (open server)
password=acell             ; default is no password (open server)

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=10MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=true                ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
username=aaa                 ; should be same as in [*_http_server] if set
password=acell               ; should be same as in [*_http_server] if set

[program:tensorboard]
command=/bin/uv run tensorboard --logdir /workspace/logs --host 0.0.0.0 --port 6006
process_name=%(program_name)s
numprocs=1
redirect_stderr=true
stdout_logfile=/tmp/tensorboard.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10

[program:jupyterlab]
command=/bin/uv run jupyter-lab
process_name=%(program_name)s
numprocs=1
redirect_stderr=true
stdout_logfile=/tmp/jupyterlab.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10

[program:sshd]
command=/usr/sbin/sshd -D
process_name=%(program_name)s
numprocs=1
redirect_stderr=true
stdout_logfile=/tmp/sshd.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10

[program:ttyd]
command=/bin/ttyd --port 13121 --writable --client-option fontFamily=monospace /bin/zsh --login
process_name=%(program_name)s
numprocs=1
redirect_stderr=true
stdout_logfile=/tmp/ttyd.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
